# vagrant-config/upstart-scripts/jiffylube-web.conf
#!upstart

description "jiffylube-web"
author "Sq1"

start on startup
stop on shutdown

respawn
respawn limit 15 5

script
        export HOME='/vagrant/jiffylube-web'

        echo $$ > /var/run/jiffylube-web.pid
        exec sudo -u vagrant NODE_ENV=local PORT=3200  coffee /vagrant/jiffylube-web/app.coffee >> /var/log/jiffylube-web.log 2>&1
end script

pre-start script
        echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting jiffylube-web"
        echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/jiffylube-web.log
end script

pre-stop script
        echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping jiffylube-web"
        echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/jiffylube-web.log
end script
